# build/supervisord.conf
# Supervisord runs as a foreground process and (re)starts all required services
# Database (pg), Web framework (Django), Image Server and Web Server (Nginx)
# Without a constant running foreground process, many container executers
# will kill the container due to inactivity.
# TODO: redirect outputs from the processes to log files
#
# To Start:
# supervisord -c build/supervisord.conf
#
# To Stop: CTRL-BREAK or Kill
# To reload config: kill -HUP
#
[supervisord]
nodaemon=true

[program:pgsql]
command=/usr/lib/postgresql/9.5/bin/postgres -D /var/lib/postgresql/9.5/main -c config_file=/etc/postgresql/9.5/main/postgresql.conf
user=postgres
stopsignal=INT

[program:iipsrv]
command=/root/iipsrv/src/iipsrv.fcgi --bind 127.0.0.1:9000 --backlog 1024
environment=LOGFILE="/tmp/iipsrv.log",VERBOSITY="1",MAX_IMAGE_CACHE_SIZE="20",JPEG_QUALITY="75",MAX_CVT="800",FILESYSTEM_PREFIX="/home/digipal/images/"
#Won't start as www-data, leave it as root
#user=www-data

[program:django]
command=/usr/bin/uwsgi --ini /home/digipal/digipal/wsgi.template.ini

[program:nginx]
# Minor issue if nginx dies, it's workers are still alive and prevent supervisor
# from restarting nginx. This is a known issue with nginx.
command=/usr/sbin/nginx -g "daemon off;"
# Will run workers as www-data, no need to force it here
