# build/supervisord.conf
# Supervisord runs as a foreground process and (re)starts all required services
# Database (pg), Web framework (Django), Image Server and Web Server (Nginx)
# Without a constant running foreground process, many container executers
# will kill the container due to inactivity.
# TODO: redirect outputs from the processes to log files
#
# To Start:
# supervisord -c build/supervisord.conf
#
# To Stop: CTRL-BREAK or Kill
# To reload config: kill -HUP
#
[supervisord]
nodaemon=true

[program:pgsql]
# TODO: remove hard-coded reference to path...
command=/usr/lib/postgresql/10/bin/postgres -D /home/digipal/digipal_project/database -c config_file=/etc/postgresql/10/main/postgresql.conf
user=postgres
stopsignal=INT
priority=100

[program:iipsrv]
command=/usr/local/bin/iipsrv.fcgi --bind 127.0.0.1:9000 --backlog 1024
environment=LOGFILE="/tmp/iipsrv.log",VERBOSITY="1",MAX_IMAGE_CACHE_SIZE="20",JPEG_QUALITY="75",MAX_CVT="800",FILESYSTEM_PREFIX="/home/digipal/digipal_project/images/"
priority=200
#Won't start as www-data, leave it as root
#user=www-data

[program:django]
command=/usr/bin/uwsgi --ini /home/digipal/build/wsgi.ini
# without this the master process still runs after exit (just reloaded content)
stopsignal = QUIT
priority=300

[program:nginx]
# Minor issue if nginx dies, it's workers are still alive and prevent supervisor
# from restarting nginx. This is a known issue with nginx.
command=/usr/sbin/nginx -g "daemon off;"
priority=400
# Will run workers as www-data, no need to force it here
